# etcd-snapshot.timer
[Unit]
Description=Run etcd-snapshot.service daily

[Timer]
OnCalendar=00:01:00 America/New_York
